# docs/nanoclaw-backup.service
[Unit]
Description=NanoClaw backup (store + db)

[Service]
Type=oneshot
WorkingDirectory=/root/nanoclaw
ExecStart=/root/nanoclaw/scripts/backup-store.sh
